# cloud-init.service
[Unit]
Description=cloud-init script
After=network.target

[Service]
ExecStart=/usr/sbin/cloud-init

[Install]
WantedBy=multi-user.target
