# cloud-init.service
[Unit]
Description=cloud-init script to configure nodes on boot

[Service]
ExecStart=/usr/sbin/cloud-init

[Install]
WantedBy=multi-user.target
